# .ebextensions/conf/upload_to_S3.conf
#!upstart

description "django-server"
author      "simona"

#start on shutdown:
start on runlevel [056]

script
        exec /emind/ebextensions/scripts/cron-jobs/log-rotate.sh
end script

